# occdriver/udev.rules
KERNEL=="snsocb[0-9]*",	MODE="0666", RUN+="/etc/snsocb/$DEVNAME"
